# backend/docker/supervisord.render.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:auth]
command=/usr/local/bin/python -m uvicorn auth.main:app --host 127.0.0.1 --port 8001
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=100
environment=PYTHONUNBUFFERED="1"

[program:book]
command=/usr/local/bin/python -m uvicorn book.main:app --host 127.0.0.1 --port 8002
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=101
environment=PYTHONUNBUFFERED="1"

[program:delivery]
command=/usr/local/bin/python -m uvicorn delivery.main:app --host 127.0.0.1 --port 8003
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=102
environment=PYTHONUNBUFFERED="1"

[program:robot]
command=/usr/local/bin/python -m uvicorn robot.main:app --host 127.0.0.1 --port 8004
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=103
environment=PYTHONUNBUFFERED="1"

[program:notification]
command=/usr/local/bin/python -m uvicorn notification.main:app --host 127.0.0.1 --port 8005
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=104
environment=PYTHONUNBUFFERED="1"

[program:celery]
; Default concurrency=1 — prefork workers each copy the app; 512MB instances OOM if left at CPU-based default (e.g. 16). Override with CELERY_CONCURRENCY on Render.
command=/bin/sh -c 'exec /usr/local/bin/python -m celery -A shared.celery_app:celery_app worker --loglevel=INFO --concurrency=${CELERY_CONCURRENCY:-1}'
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=200
environment=PYTHONUNBUFFERED="1"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=400
